# deploy/systemd/b3sr-ingest.service
[Unit]
Description=B3SR ingest from S3 to VictoriaMetrics
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/ingest/current
EnvironmentFile=-/etc/ingest/ingest.env
ExecStart=/home/ubuntu/ingest/current/deploy/bin/run_b3sr_ingest.sh
StandardOutput=append:/home/ubuntu/ingest/logs/b3sr.log
StandardError=append:/home/ubuntu/ingest/logs/b3sr.log
Nice=10

[Install]
WantedBy=default.target
