# deploy/systemd/b3sr-ingest.timer
[Unit]
Description=Run B3SR ingest every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=b3sr-ingest.service

[Install]
WantedBy=timers.target
